# ezgoo/supervisor/ezgoo.conf
[program:ezgoo]
command=/var/www/ezgoo/ezgoo -dir /var/www/ezgoo/dist
autorestart=true
user=nobody
